# scripts/install/docker-volume-vsphere.conf
description "Docker Volume Driver for vSphere"

start on starting docker
stop on stopped docker RESULT=ok

respawn

kill timeout 20

exec /usr/local/bin/$UPSTART_JOB
